# com.aliir74.x-shield.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aliir74.x-shield</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/aliirani/.local/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-u</string>
        <string>-m</string>
        <string>src.shield</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/aliirani/Downloads/Coding/x-shield</string>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/aliirani/Downloads/Coding/x-shield/cron.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/aliirani/Downloads/Coding/x-shield/cron.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
